# setup/mysqld-supervisor.conf
[program:mysqld]
command=/bin/bash -c "su mysql -s /bin/bash -c 'mysqld'"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
